# gnome/nostr-goa-overlay/provision/systemd_user/nostrfs.service
[Unit]
Description=Nostr Filesystem (user)
After=nostr-router.service
Requires=nostr-router.service

[Service]
Type=simple
EnvironmentFile=-%h/.config/nostr-goa-overlay/env
ExecStart=/usr/bin/env bash -lc '
  BIN="${NOSTR_FS_BIN:-}"
  if [ -n "$BIN" ]; then exec "$BIN" %h/Nostr; fi
  if [ -x %h/.local/bin/nostrfs ]; then exec %h/.local/bin/nostrfs %h/Nostr; fi
  exec nostrfs %h/Nostr'
Restart=on-failure
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077

[Install]
WantedBy=default.target
