# gnome/nostr-homed/systemd/nostrfs@.service
[Unit]
Description=NostrFS FUSE mount for user %i
After=network-online.target
Requires=network-online.target
BindsTo=nostr-homectl.service

[Service]
Type=simple
Environment=HOMED_NAMESPACE=personal
Environment=SECRETS_NAMESPACE=personal
Environment=BLOSSOM_BASE_URL=https://blossom.example.org
Environment=RELAYS_DEFAULT=wss://relay.damus.io,wss://nostr.wine
ExecStart=/usr/bin/nostrfs --writeback --namespace=${HOMED_NAMESPACE} /home/%i
Restart=on-failure
RestartSec=2s
# Sandboxing
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ProtectHome=read-only
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictNamespaces=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077
DeviceAllow=/dev/fuse rwm
CapabilityBoundingSet=CAP_SYS_ADMIN
AmbientCapabilities=CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
